# timers/auto-btrfs-snapshot-weekly.service
[Unit]
Description=Take read-only Btrfs weekly snapshots of media

[Service]
Type=oneshot
ExecStart=/usr/local/bin/auto-btrfs-snapshot-weekly.sh
